# examples/journald-notify.service
[Unit]
Description=journald-notify
After=network.target
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/bin/journald-notify -c /etc/journald-notify.json run
User=journald-notify

[Install]
WantedBy=multi-user.target
